# themes/luci-theme-supernet/root/etc/uci-defaults/30_luci-theme-supernet
#!/bin/sh
uci batch <<-EOF
	set luci.themes.SuperNet=/luci-static/supernet
	set luci.main.mediaurlbase=/luci-static/supernet
	commit luci
EOF
exit 0 
# ↑↑↑ DO NOT EDIT THIS LINE
